# tools/llvmlinux/targets/malta/initramfs-overlay/etc/inittab
::sysinit:/etc/init.d/rcS
::sysinit:/etc/init.d/autopoweroff
::askfirst:-/bin/sh
